# rootfs/etc/services.d/synapse/run
#!/usr/bin/with-contenv bash
set -Eeuo pipefail
python -m synapse.app.homeserver \
    -c /data/homeserver.yaml \
    --data-directory /data \
    --config-directory /data \
