# rootfs/etc/services.d/synapse/run
#!/usr/bin/with-contenv bash
set -Eeuo pipefail

exec /start.py
